# crontab
*/1 * * * * root /etc/cron.custom/keep-alive.sh
